# templates/downloads-delete.timer
[Unit]
Description=Periodically clean Downloads

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
