# templates/downloads-delete.service
[Unit]
Description=Delete old files in Downloads

[Service]
Type=simple
ExecStart=/usr/bin/find $HOME/Downloads -mtime +30 -delete
